# systemd/user/canbusd.service
[Unit]
Description=open-mmi CAN Bus Daemon
Documentation=https://github.com/Sheepdog-97/open-mmi
After=graphical-session.target

[Service]
Type=simple

WorkingDirectory=/opt/open-mmi

# Python logging buffering
Environment="PYTHONUNBUFFERED=1"

# Logging level (DEBUG, INFO, WARNING, ERROR)
Environment="OPEN_MMI_LOG_LEVEL=INFO"

# Vehicle profile (default: seat-leon-1p-pq35)
Environment="OPEN_MMI_VEHICLE=seat-leon-1p-pq35"

# Bindings file (default: default)
Environment="OPEN_MMI_BINDINGS=default"

# CAN runtime selection
# Defaults are equivalent to:
#   OPEN_MMI_CAN_BUS=comfort
#   OPEN_MMI_CAN_INTERFACE=can0
#
# Override locally with:
#   ./scripts/manage.sh config edit-service
#
# Example override:
#   Environment="OPEN_MMI_CAN_BUS=comfort"
#   Environment="OPEN_MMI_CAN_INTERFACE=vcan0"

# Execution
ExecStart=/opt/open-mmi/venv/bin/python -m canbusd.core

# Restart policy
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=canbusd

[Install]
WantedBy=default.target
